# recipes-qt/qt4/qt4-x11-free-systemd/qtdemo.service
[Unit]
Description=start X and qtdemo

[Service]
Type=forking
ExecStart=/usr/bin/qtdemo-init start
ExecStop=/usr/bin/qtdemo-init stop

[Install]
WantedBy=graphical.target
